# systemd/ai-relay-service.service
[Unit]
Description=AI Relay Server v2
After=network-online.target

[Service]
Type=simple
ExecStart=%h/projects/ai-relay-service/.venv/bin/relay-server server --port 8788
WorkingDirectory=%h/projects/ai-relay-service
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=RELAY_DB_PATH=%h/.relay/server.db
# REQUIRED: Set RELAY_SESSION_SECRET or create ~/.relay/config.yaml with
# session_secret. Without it the server refuses to start.
# Environment=RELAY_SESSION_SECRET=your-32-char-secret-here

[Install]
WantedBy=default.target
